# edge-node/systemd/nightforge-shield.service
[Unit]
Description=NightForge Shield — Suricata + Cowrie threat scoring and auto-blocking
After=network-online.target nightforge-cowrie.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/bin/bash /home/foreverlx/scripts/suricata/suricata-shield.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
