# edge-node/containers/caddy.container
[Unit]
Description=Caddy TLS Reverse Proxy
After=network-online.target vaultwarden.service homepage.service
Wants=network-online.target

[Container]
Image=docker.io/library/caddy:latest
ContainerName=caddy
Network=host
Volume=%h/caddy/Caddyfile:/etc/caddy/Caddyfile:ro
Volume=%h/caddy/data:/data:z
Volume=%h/caddy/config:/config:z
AutoUpdate=registry
HealthCmd=curl -sf http://localhost:2019/config/
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3

[Service]
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
