# edge-node/containers/vaultwarden.container
[Unit]
Description=Vaultwarden Password Manager
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/vaultwarden/server:latest
ContainerName=vaultwarden
PublishPort=8081:80
Volume=%h/vaultwarden:/data:z
Environment=WEBSOCKET_ENABLED=true
AutoUpdate=registry
HealthCmd=curl -sf http://localhost:80/alive
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3

[Service]
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
